# grocy/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Community Hass.io Add-ons: Grocy
# Runs the Nginx daemon
# ==============================================================================

# Wait for PHP-FPM to become available
bashio::net.wait_for 9001

bashio::log.info "Starting NGinx...."

exec nginx
