# grocy/rootfs/etc/s6-overlay/s6-rc.d/nginx/run
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: Grocy
# Runs the Nginx daemon
# ==============================================================================

# Wait for PHP-FPM to become available
bashio::net.wait_for 9002

bashio::log.info "Starting NGinx...."

exec nginx
